# persistent.sh
#!/bin/bash
# chkconfig: 2345 20 80
# description: 

# Source function library.

# Check for updates and install them
# Renew certs for nginx, vsftpd etc.